# SensorMonitor.service
[Unit]
Description=SensorMonitor Service (Pi Zero Sensor Node)
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /home/pi/SensorMonitor.py
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
